# .config/local/supervisor/supervisor-uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/app/bin/uwsgi -i /srv/app/.config/local/uwsgi/app.ini

