# .config/test/supervisor/supervisor-uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/app/bin/uwsgi -i /srv/app/.config/test/uwsgi/app.ini

