# .config/deploy/supervisor/supervisor-rabbitmq.conf
[program:rabbitmq]
command = rabbitmq-server
