# .config/deploy/supervisor/supervisor-uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/app/bin/uwsgi -i /srv/app/.config/deploy/uwsgi/app.ini

