# .config/deploy/supervisor/supervisor-celery.conf
[program:celery]
command = /root/.pyenv/versions/app/bin/celery -A config worker --loglevel=info
directory = /srv/app/petweb